# loadavg/loadavg.service
[Unit]
Description=loadavg

[Service]
ExecStart=/loadavg

[Install]
WantedBy=multi-user.target